# 425500_Cloud_Computing_Source_Code/services/spark_python_job.service
[Unit]
Description=Script to run the ingestion python script

[Service]
ExecStart=/usr/local/bin/start_spark_job.sh

[Install]
WantedBy=multi-user.target